# pg-tomcat/dockerbuild/pgconfig/postgresql-service.conf
author "gsmiro"
start on startup
stop on runlevel [016]
exec $PGBIN/bin/pg_ctl -D $PGCLUSTER -l /var/log/metaway_postgresql.log
